# deploy/finadvisorbot.service
[Unit]
Description=FinAdvisor Telegram Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/FinAdvisor_AI_bot
# Активация venv через bash — не зависит от имени бинарника (python/python3)
# Если venv в другой папке (например .venv), замените venv на .venv ниже
ExecStart=/bin/bash -c '. /root/FinAdvisor_AI_bot/venv/bin/activate && exec python bot.py'
Restart=always
RestartSec=5
EnvironmentFile=/root/FinAdvisor_AI_bot/.env

[Install]
WantedBy=multi-user.target
