# deploy/finadvisor-api.service
[Unit]
Description=FinAdvisor API (FastAPI + React)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/FinAdvisor_AI_bot
ExecStart=/root/FinAdvisor_AI_bot/venv/bin/python -m uvicorn api:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5
EnvironmentFile=/root/FinAdvisor_AI_bot/.env

[Install]
WantedBy=multi-user.target
