# deploy/finadvisorbot.service
[Unit]
Description=FinAdvisor Telegram Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/FinAdvisor_AI_bot
# Скрипт сам находит venv или .venv в каталоге проекта
ExecStart=/root/FinAdvisor_AI_bot/scripts/run_bot_venv.sh
Restart=always
RestartSec=5
EnvironmentFile=/root/FinAdvisor_AI_bot/.env

[Install]
WantedBy=multi-user.target
